# packages/java/rinfo-sesame-http/conf/tomcat/tomcat
#!/bin/sh
#
# Startup script for Tomcat Servlet Engine
#
# chkconfig: 345 86 14
# description: Tomcat Servlet Engine
# processname: tomcat
# pidfile: /usr/share/tomcat6/bin/tomcat.pid
#

# Source function library.
. /etc/rc.d/init.d/functions

# User under which tomcat will run
TOMCAT_USER=tomcat
CATALINA_HOME=/usr/share/tomcat6

RETVAL=0

# start, debug, stop, and status functions
start() {
    # Start Tomcat in normal mode
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
    if [ $SHUTDOWN_PORT -ne 0 ]; then
        echo "Tomcat already started"
    else
        echo "Starting tomcat..."
        chown -R $TOMCAT_USER:$TOMCAT_USER /usr/share/tomcat6/*
        chown -R $TOMCAT_USER:$TOMCAT_USER /home/tomcat/*
        chown -R $TOMCAT_USER:$TOMCAT_USER /opt/_workapps/rinfo/aduna
        su -l $TOMCAT_USER -c '/usr/share/tomcat6/bin/catalina.sh start'
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        while [ $SHUTDOWN_PORT -eq 0 ]; do
            echo "Waiting for Tomcat to start" 
            sleep 1
        SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        done    
        RETVAL=$?
        echo "Tomcat started in normal mode"
        [ $RETVAL=0 ] && touch /var/lock/subsys/tomcat
    fi
}

debug() {
    # Start Tomcat in debug mode
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
    if [ $SHUTDOWN_PORT -ne 0 ]; then
        echo "Tomcat already started"
    else
        echo "Starting tomcat in debug mode..."
        chown -R $TOMCAT_USER:$TOMCAT_USER /usr/share/tomcat6/*
        chown -R $TOMCAT_USER:$TOMCAT_USER /home/tomcat/*
        su -l $TOMCAT_USER -c '/usr/share/tomcat6/bin/catalina.sh jpda start'
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        while [ $SHUTDOWN_PORT -eq 0 ]; do
            sleep 1
            SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        done    
        RETVAL=$?
        echo "Tomcat started in debug mode"
        [ $RETVAL=0 ] && touch /var/lock/subsys/tomcat
    fi
}

stop() {
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
    if [ $SHUTDOWN_PORT -eq 0 ]; then
        echo "Tomcat already stopped"
    else
        echo "Stopping tomcat..."
        su -l $TOMCAT_USER -c '/usr/share/tomcat6/bin/catalina.sh stop'
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        while [ $SHUTDOWN_PORT -ne 0 ]; do
            sleep 1
            SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
        done
    RETVAL=$?
        echo "Tomcat stopped"
        [ $RETVAL=0 ] && rm -f /var/lock/subsys/tomcat /usr/share/tomcat6/bin/tomcat.pid
    fi
}

status() {
    SHUTDOWN_PORT=`netstat -vatn|grep LISTEN|grep 9915|wc -l`
    if [ $SHUTDOWN_PORT -eq 0 ]; then
        echo "Tomcat stopped"
    else
        MODE="normal"
        JPDA_PORT=`netstat -vatn|grep LISTEN|grep 8000|wc -l`
        if [ $JPDA_PORT -ne 0 ]; then
            MODE="debug"
        fi
    echo "Tomcat running in $MODE mode"
    fi
}

case "$1" in
  start)
        start
        ;;
  debug)
        debug
        ;;
  stop)
        stop
        ;;
  restart)
        stop
        start
        ;;
  redebug)
        stop
        debug
        ;;
  status)
    status
    ;;
  *)
    echo "Usage: $0 {start|debug|stop|restart|redebug|status}"
    exit 1
esac

exit $RETVAL
